# launchd/com.EXAMPLE.assistant.maintenance.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.EXAMPLE.assistant.maintenance</string>
    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_DIR__/.venv/bin/python3</string>
        <string>__INSTALL_DIR__/maintenance.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>
    <key>RunAtLoad</key>
    <false/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__INSTALL_DIR__/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/maintenance.log</string>
    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/maintenance.log</string>
</dict>
</plist>
