# launchd/com.EXAMPLE.assistant.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.EXAMPLE.assistant.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_DIR__/.venv/bin/python3</string>
        <string>__INSTALL_DIR__/moleck.py</string>
        <string>--serve</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>DALIA_API_HOST</key>
        <string>127.0.0.1</string>
        <key>PATH</key>
        <string>__INSTALL_DIR__/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/server.log</string>
    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/server.log</string>
</dict>
</plist>
